# crontab
1 0 * * * /schedule.sh >> /var/log/cron.log 2>&1
